# src/etc/init/ipc-shutdown.conf
# ipc-shutdown - Shutdown the ipcsystools
#
# This task is run on system shutdown

description	"Shutdown IPC enviroment"

start on runlevel [06]

task

exec /usr/sbin/ipc-shutdown
