# src/etc/init/ipc-shutdown.conf
# ipc-shutdown - Shutdown the ICP
#
# This task is run on shutdown the IPC

description	"Shutdown IPC enviroment"

start on runlevel [06]

task

exec /etc/init.d/ipc-shutdown
