# src/etc/cron.d/ipcsystools
# Installed by ipcsystools

*/2 * * * * root /usr/share/ipcsystools/ipc-cronjobs.sh
