# src/etc/if-pre-up.d/01-3g-led-blink
#!/bin/bash

# Basis Bibliothek für die IPCs
. /usr/share/ipcsystools/ipclib.inc

if check_ifupdown_on_mobile_connection ; then 
	syslogger "debug" "interface $IFACE : 3G led blink"
	/usr/bin/ipc-set-led 3g timer
else
	syslogger "debug" "interface $IFACE : no action"
fi

exit 0


