# systemd/ipcsystools-monitor.service
[Unit]
Description=Supervise and manage network and VPN connections for routing

[Service]
Type=oneshot
ExecStart=/usr/bin/ipc-monitor
